# docker-ovs-plugin/rootfs/var/lib/dpkg/info/uuid-runtime.list
/.
/etc
/etc/init.d
/etc/init.d/uuidd
/usr
/usr/bin
/usr/bin/uuidgen
/usr/bin/uuidparse
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/uuidd.service
/usr/lib/systemd/system/uuidd.socket
/usr/lib/sysusers.d
/usr/lib/sysusers.d/uuidd-sysusers.conf
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/uuidd-tmpfiles.conf
/usr/sbin
/usr/sbin/uuidd
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/uuidd
/usr/share/bash-completion/completions/uuidgen
/usr/share/bash-completion/completions/uuidparse
/usr/share/doc
/usr/share/doc/uuid-runtime
/usr/share/doc/uuid-runtime/NEWS.Debian.gz
/usr/share/doc/uuid-runtime/changelog.Debian.gz
/usr/share/doc/uuid-runtime/changelog.gz
/usr/share/doc/uuid-runtime/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/uuidgen.1.gz
/usr/share/man/man1/uuidparse.1.gz
/usr/share/man/man8
/usr/share/man/man8/uuidd.8.gz
